# myhubot/upstart/myhubot.conf
description "MYHUBOT"
author "ME <leonmakster@gmail.com>"

start on runlevel [2345]
stop on runlevel [016]

respawn

setuid vagrant

env HOME='/home/vagrant'

chdir /vagrant/myhubot/

script
  export PATH="node_modules/.bin:node_modules/hubot/node_modules/.bin:$PATH"
  export HUBOT_HIPCHAT_JID=503991_3253123@chat.hipchat.com
  export HUBOT_HIPCHAT_PASSWORD=lmakerlz

  exec node_modules/.bin/hubot --adapter hipchat --name "myhubot"

end script
